# generate-index.timer
[Unit]
Description=Runs generate-index.service everyday at 05:00 PST

[Timer]
OnCalendar=*-*-* 05:00:00 America/Vancouver
Persistent=true
